# dist/rtlamr.service
[Unit]
Description=Service to capture power usage from rtl_tcp
Documentation="https://github.com/bemasher/rtlamr"
After=rtl-tcp.service
Wants=rtl-tcp.service

[Service]
ExecStart=/opt/powermeter/start.sh
Type=Simple
Restart=on-failure
User=root

[Install]
WantedBy=multi-user.target